# scripts/systemd/rhumb-agent.service.tmpl
[Unit]
Description=Rhumb agent host (Claude Code session API)
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
User=@RUN_USER@
WorkingDirectory=@REPO_DIR@/agent-host
EnvironmentFile=/etc/rhumb/rhumb.env
ExecStart=@NODE_BIN@ dist/index.js
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
